# recipes-core/microdev-startup-service/microdev-startup-service/microdev.service
[Unit]
Description=Start LAE Script

[Service]
ExecStart=/usr/bin/startUpMicrodev.sh

[Install]
WantedBy=multi-user.target
